# src/frontend/mtcamlx-frontend.socket
[Unit]
Description=FAUST CTF helloworld service

[Socket]
ListenStream=7777
Accept=true

[Install]
WantedBy=sockets.target

